# system/etc/udev/rules.d/70-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/jc/.scripts/monitor-switch.sh"
